# service/discord-all-in-backup.timer
[Unit]
Description=Hourly snapshot of the discord-all-in database

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
